# rootfs/etc/s6/nginx/finish
#!/bin/bash
/usr/bin/s6-svscanctl -t /etc/s6
